# docker/logstash.ini
[program:logstash]
command=/usr/share/logstash/bin/logstash --path.settings=/etc/logstash
directory=/usr/share/logstash
environment=JAVACMD=/usr/bin/java
redirect_stderr=true
stopsignal=QUIT
stopwaitsecs=30
user=root