# docker/rabbitmq.ini
[program:rabbitmq-server]
command=/usr/sbin/rabbitmq-server
directory=/opt
redirect_stderr=true
stopsignal=QUIT
stopwaitsecs=30
user=root