# docker/rabbitmq.ini
[program:rabbitmq-server]
command=/usr/sbin/rabbitmq-server
directory=/opt
redirect_stderr=true
stopsignal=INT
stopwaitsecs=5
user=root
